# data-collection-master/Project-App/src/output/temp/lib/systemd/system/app-daemon.service
[Unit]
Description=app-daemon service

[Service]
ExecStart=/usr/bin/app-daemon
ExecReload=/bin/kill -HUP 
KillMode=process
Restart=always

[Install]
WantedBy=multi-user.target
